# roles/systemd/files/ansible-pull@.timer
[Timer]
OnCalendar=hourly
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=multi-user.target

